# rootfs/etc/s6-overlay/s6-rc.d/otbr-agent/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# OpenThread BorderRouter start script
# ==============================================================================

# shellcheck disable=SC1091
. /etc/s6-overlay/scripts/otbr-agent-common

declare backbone_if
declare thread_if
# shellcheck disable=SC2034
declare device
# shellcheck disable=SC2034
declare baudrate
# shellcheck disable=SC2034
declare flow_control
declare otbr_rest_listen
declare otbr_rest_listen_port

# TODO: make configurable?
backbone_if="wpan0"
thread_if="wpan0"

# shellcheck disable=SC2015
mkdir -p /data/thread && ln -sft /var/lib /data/thread || bashio::exit.nok "Could not create directory /data/thread to store Thread data."

# TODO
# Firewall has to be configured on the host, not inside the container
# Create the ipsets anyways, to prevent errors from otbr-agent
ipset create -exist otbr-ingress-deny-src hash:net family inet6
ipset create -exist otbr-ingress-deny-src-swap hash:net family inet6
ipset create -exist otbr-ingress-allow-dst hash:net family inet6
ipset create -exist otbr-ingress-allow-dst-swap hash:net family inet6

# TODO: make port configurable (otbr-agent refuses to start if the port is changed)
otbr_rest_listen="::"
otbr_rest_listen_port="8081"

# Store REST API listen information for check script
echo "${otbr_rest_listen}" > /tmp/otbr-agent-rest-api
echo "${otbr_rest_listen_port}" >> /tmp/otbr-agent-rest-api

bashio::log.info "Starting otbr-agent on port ${otbr_rest_listen_port}..."
# shellcheck disable=SC2086
exec s6-notifyoncheck -d -s 300 -w 300 -n 0 \
    "/usr/sbin/otbr-agent" -I ${thread_if} -B "${backbone_if}" \
        --rest-listen-address "${otbr_rest_listen}" \
        -d${OTBR_LOG_LEVEL} -v \
        "spinel+cpc://cpcd_0?iid=2&iid-list=0"
