# rootfs/etc/s6-overlay/s6-rc.d/zigbeed-tcp/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start EmberZNet Zigbee Daemon
# ==============================================================================
# shellcheck disable=SC2015

# TODO
ezsp_listen_port="9999"
#ezsp_listen_port="$(bashio::addon.port 9999)"

bashio::log.info "Starting zigbeed-tcp..."
exec /usr/bin/socat TCP-LISTEN:${ezsp_listen_port},reuseaddr,fork FILE:/tmp/ttyZigbee,raw,echo=0,b115200
