# rootfs/etc/s6-overlay/s6-rc.d/zigbeed-tcp/run
#!/usr/bin/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start EmberZNet Zigbee Daemon
# ==============================================================================
# shellcheck disable=SC2015

bashio::log.info "Starting zigbeed-tcp on port ${ZIGBEED_PORT}..."
exec /usr/bin/socat TCP-LISTEN:${ZIGBEED_PORT},reuseaddr,fork FILE:/tmp/ttyZigbee,raw,echo=0,b115200
